# prod/files/zuluru-sync.timer
[Unit]
Description=Run Zuluru Sync Daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
